# systemd/snap-sync@.service
[Unit]
Description=snap-sync (UUID: %i)

[Install]
WantedBy=multi-user.target

[Service]
Type=simple

ExecStart=/usr/bin/snap-sync-wrapper %i
